# config/includes.binary/etc/rc.local
#!/bin/sh -e
gradm2 -F -L /etc/grsec/learning.logs &
exit 0
