# config/includes.binary/etc/rc.local
#!/bin/sh -e
adduser \ grsec-tpe &
adduser \ fuse &
adduser \ docker &
exit 0
